# airootfs/etc/systemd/system/autoinstall.service
[Unit]
Description=Eadxm Automated Arch Architect
# Wait for the boot sequence to settle
After=systemd-user-sessions.service getty@tty1.service
# Prevent the default Arch login prompt from fighting us for the keyboard
Conflicts=getty@tty1.service

[Service]
# 'idle' ensures systemd waits until all other boot text stops before launching
Type=idle
ExecStart=/usr/local/bin/auto-install.sh

# Route the script to the physical monitor and keyboard
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1

# Prevent systemd from killing the script if it takes longer than 90 seconds
TimeoutStartSec=0
User=root

[Install]
WantedBy=multi-user.target
